# deploy/systemd/lobster-bot.service
[Unit]
Description=lobster-bot Telegram Assistant
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/lobster-bot
EnvironmentFile=/root/lobster-bot/.env
ExecStart=/usr/bin/python3 run.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
